# deploy/systemd/octo-api-front.service
[Unit]
Description=Octo Agent Test API Front
After=network.target octo-api.service

[Service]
Type=simple
User=octo
Group=octo
WorkingDirectory=/opt/octo-agent/current/frontend
EnvironmentFile=/etc/octo-agent/test/api-front.env
Environment=PATH=/opt/octo-agent/runtime/node/current/bin:/opt/octo-agent/runtime/go/current/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/octo-agent/current/frontend/node_modules/.bin/next start -p 4200
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
